# pkg/postinst
#!/bin/sh
# Enable the systemd service
systemctl enable kbled.service
# Start the systemd service
systemctl start kbled.service
# Exit successfully
exit 0
